# test.ini
[program:ic-publisher-track]
command="php core.php Consumer --channels=opensclicks/%engine% GenericJobProcessor"
process_name="track-engine-%engine%-%(process_num)02d"
numprocs=2
engine=0..2
directory=/var/www/icp
user=icp
stopsignal=INT
